# supervisor.conf
[supervisord]
logfile=/proc/self/fd/1
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:rpcbind]
command=/sbin/rpcbind -wdh 0.0.0.0 -f
autostart=true
autorestart=true
pidfile=/var/run/supervisor-rpcbind.pid
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:rpcnfsd]
command=/usr/sbin/rpc.nfsd -d -V 4.1 -N 2 -N 3 -N 4.0
autostart=true
autorestart=false
pidfile=/var/run/supervisor-rpcnfsd.pid
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:exportfs]
command=/usr/sbin/exportfs -rv
autostart=true
autorestart=false
pidfile=/var/run/supervisor-exportfs.pid
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0

[program:rpcmountd]
command=/usr/sbin/rpc.mountd --debug all -V 4.1 -N 2 -N 3 -F
autostart=true
autorestart=true
pidfile=/var/run/supervisor-rpcmountd.pid
stdout_logfile=/proc/self/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/proc/self/fd/2
stderr_logfile_maxbytes=0
